# q.conf
[program:q]
command=python3 /home/curator2/q/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/q/test.err.log
stdout_logfile=/home/curator2/q/test.out.log